# .seed/etc/cron.d/ipv6_addr_notify
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * *   root	/usr/bin/ipv6_addr_notify.sh >/dev/null
