# .seed/etc/cron.d/vpn_check
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/1 * * * *   root	/usr/bin/vpn_check.sh >/dev/null
